# files/ttyS1.conf
# This file is being maintained by Puppet.
# DO NOT EDIT
#
# ttyS1 - agetty
#
# This service maintains a agetty on ttyS1.

stop on runlevel [S016]
start on runlevel [23]

respawn
exec /sbin/agetty -h -L -w /dev/ttyS1 115200 vt102
